# ops/systemd/persian-tools-ops-snapshot-production.timer
[Unit]
Description=Run Persian Tools production ops snapshot every 6 minutes

[Timer]
Unit=persian-tools-ops-snapshot-production.service
Persistent=true
OnBootSec=5min
OnUnitActiveSec=6min

[Install]
WantedBy=timers.target
